# redhat/paas/demo/cygwin/supervisor/conf.d/5022.ini
[program:5022]
command=/usr/bin/autossh -M 0 -N -L 5022:apps1.paas.com:22 -i /Users/wzh/.ssh/id_rsa -p 22 root@172.16.20.143
autorestart=true
