# a6test/hadoop/conf/hive/supervisord.conf
[program:metastore]
command=/opt/hive/bin/hive --service metastore

[program:hiveserver2]
command=/opt/hive/bin/hive --service hiveserver2